# deploy/cogmemory-test.service.d/socket-perms.conf
[Service]
# Make the test socket world-accessible so agents running under any uid can connect.
ExecStartPost=/bin/chmod 777 %h/.local/share/cogmemory-test/cog-memory-test.sock
